# zabbix-agent2.4/web/supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D

[program:tomcat]
command = /usr/local/tomcat/bin/catalina.sh run
stopasgroup = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:zabbix_agentd]
command=/usr/local/zabbix_agent/sbin/zabbix_agentd
user=zabbix
startsecs=0
autorestart=false
redirect_stderr = true
